# cron/BackupDongDanTasks.cron
#!/usr/bin/env bash

# backup dongDong in cron
# */15 * * * * /path/to/BackupDongDanTasks.cron >/dev/null 2>&1

if type greadlink > /dev/null 2>&1; then
  SHELL_FOLDER=$(dirname $(greadlink -f "$0"))
else
  SHELL_FOLDER=$(dirname $(readlink -f "$0"))
fi

cd $SHELL_FOLDER
PHP_BIN=$(which php)
$PHP_BIN ../public/index.php BackupDongDan  
$PHP_BIN ../public/index.php BackupDongDanComments "take=20"
$PHP_BIN ../public/index.php BackupDongDanImages "take=20"

./UpdateDongDanOldImages.cron

exit 0

