# cron/PubHackerNewsViaBaiduFy.cron
#!/usr/bin/env bash

# update old 50x50 portaits to 200x200
# */15 * * * * /path/to/UpdateDongDanOldImages.cron >/dev/null 2>&1

if type greadlink > /dev/null 2>&1; then
  SHELL_FOLDER=$(dirname $(greadlink -f "$0"))
else
  SHELL_FOLDER=$(dirname $(readlink -f "$0"))
fi
cd $SHELL_FOLDER

PHP_BIN=$(which php)
$PHP_BIN ../public/index.php PubHackerNews "t=baidu"

exit 0

